# files/etc/supervisor.d/kresd-gc.ini
[program:kresd-gc]
command=kres-cache-gc -c %(ENV_CACHEPATH)s -d 5000 -u 50
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startsecs=5
