# files/etc/supervisor.d/kresd-gc.ini
[program:kresd]
command=kres-cache-gc -c . -d 500
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startsecs=5