# files/etc/supervisor.d/kresd.ini
[program:kresd]
command=kresd -q -c /etc/knot-resolver/kresd.conf
process_name=%(program_name)s_%(process_num)s
environment=WORKER=%(process_num)s
numprocs=%(ENV_NUM_CPUS)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startsecs=60
